# roles/static-panda/templates/static_panda_upstart.conf
start on runlevel [2345]
stop on runlevel [016]

chdir /home/user
respawn
exec python /home/user/static-panda.py


